# cron/crontab
* * * * * /scripts/app_uptime.sh >> /var/log/cron.log 2>&1